# backup.cron
#每周六13点进行备份
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
#MAILTO=petrus
HOME=/home/petrus
0 15 * * 6 export DISPLAY=:0 ; bash /home/petrus/Backup/push2github.sh && notify-send 'System configration backup success.' || notify-send -u critical 'System configration backup failure.'
