# backup.cron
#每周六13点进行备份
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
#MAILTO=petrus
HOME=/home/petrus
#0 13 * * 6 bash /home/petrus/Backup/backup_config.sh
47 20 * * * bash /home/petrus/Backup/push2github.sh && notify-send 'Backup configration success.'
